# .build/aria2.supervisor.conf
[supervisord]
nodaemon=false

[program:background_process]
command=aria2c --enable-rpc --save-session /cloudreve/data
autostart=true
autorestart=true